# deploy/com.portfolio.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Streamlit dashboard server – starts at login and stays alive.
     Install: cp com.portfolio.dashboard.plist ~/Library/LaunchAgents/
              launchctl load -w ~/Library/LaunchAgents/com.portfolio.dashboard.plist
     Access:  http://<mac-mini>.local:8501
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.portfolio.dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>
            source ~/.zshrc 2>/dev/null;
            cd /Users/mark/dev/Portfolio;
            .venv/bin/streamlit run src/portfolio/dashboard/app.py
                --server.address 0.0.0.0
                --server.port 8501
                --server.headless true
        </string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mark/dev/Portfolio</string>

    <!-- Start at login, restart if it crashes -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/mark/data/portfolio/logs/dashboard.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mark/data/portfolio/logs/dashboard_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/mark/.pyenv/shims:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/mark/dev/Portfolio/src</string>
    </dict>
</dict>
</plist>
